# deploy/com.anchit.bhai-sunn-stt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.anchit.bhai-sunn-stt</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/anchitsom/services/bhai-sunn-stt/start.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/anchitsom/agent-brain/projects/experiments/project-bhai-sunn_asp</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HF_HOME</key>
    <string>/Users/anchitsom/.cache/huggingface</string>
    <key>HF_HUB_OFFLINE</key>
    <string>1</string>
    <key>TRANSFORMERS_OFFLINE</key>
    <string>1</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/anchitsom/services/bhai-sunn-stt/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/anchitsom/services/bhai-sunn-stt/stderr.log</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
